# Goal Finder/Unused Extra Scripts/Test and Example Scripts/atest.conf
start on startup
task
exec /home/pi/Desktop/ImageProcessing/contourfinder.py 2>/home/pi/Desktop/ImageProcessing/aimbotlog
